# sync-dir.service
[Unit]
Description=Synchronounce directories service

[Service]
Type=simple
ExecStart=~/bin/sync-dir watch
ExecStop=~/bin/sync-dir stop
